# monitor/day-trader-failsafe.service
[Unit]
Description=Day Trader Force-Close Failsafe
After=network.target

[Service]
Type=oneshot
User=harry
WorkingDirectory=/home/harry/day-trader
ExecStart=/home/harry/day-trader/venv/bin/python3 force_close_failsafe.py
EnvironmentFile=/home/harry/day-trader/.env
